# com.spinning.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.spinning.dashboard</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd ~/Dashboard && source venv/bin/activate && streamlit run app/dashboard.py --server.port 8501 --server.address 0.0.0.0 --server.headless true</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/dashboard.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/dashboard_error.log</string>
    
    <key>WorkingDirectory</key>
    <string>~/Dashboard</string>
</dict>
</plist>
